# system/systemd/backup.timer
[Unit]
Description=Run dotfiles backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
