# system/systemd/cleanup.timer
[Unit]
Description=Weekly cleanup job

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
